# airootfs/etc/init.d/sddm
#!/sbin/openrc-run
# FenrirOS - SDDM OpenRC service

name="sddm"
description="Simple Desktop Display Manager"
command="/usr/bin/sddm"
command_background=true
pidfile="/run/sddm.pid"
command_args=""

depend() {
    need dbus elogind
    use xfs
    after openrc-boot modules
    provide display-manager
}

start_pre() {
    # Ensure elogind is running before SDDM starts
    if ! rc-service elogind status > /dev/null 2>&1; then
        eerror "elogind must be running for SDDM to work correctly"
        return 1
    fi
}

start() {
    ebegin "Starting SDDM"
    start-stop-daemon --start \
        --pidfile "${pidfile}" \
        --make-pidfile \
        --background \
        --exec "${command}"
    eend $?
}

stop() {
    ebegin "Stopping SDDM"
    start-stop-daemon --stop \
        --pidfile "${pidfile}" \
        --retry TERM/5/KILL/5
    eend $?
}
